# opensand-deploy/src/backend/opensand-conf.socket
[Unit]
Description=Socket for uWSGI app OpenSAND

[Socket]
ListenStream=/run/opensand/opensand-conf.socket
SocketUser=opensand
SocketGroup=opensand
SocketMode=0660

[Install]
WantedBy=sockets.target
